# root/etc/s6-overlay/s6-rc.d/service-novnc/run
#!/command/with-contenv bash
# shellcheck shell=bash

umask "${UMASK}"

exec s6-setuidgid hotio "${APP_DIR}/utils/novnc_proxy" --vnc "${VNC_SERVER}:${VNC_PORT}" --listen "0.0.0.0:${WEBUI_PORTS%%/*}"
